# scripts/nightly/com.natemccoy.nightly-rust-clean-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.natemccoy.nightly-rust-clean-build</string>
    <!-- Fires every StartInterval seconds. The wrapper script checks HID
         idle time (via ioreg) and runs the nightly build whenever the user
         has been away from keyboard for >= 1 hour AND no previous run is
         still in flight (pgrep guard against the build script's path). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/natemccoy/.claude/scripts/nightly/nightly-trigger.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/nightly-rust-clean-build-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/nightly-rust-clean-build-stderr.log</string>
</dict>
</plist>
