# scripts/nightly/com.natemccoy.nightly-rust-clean-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.natemccoy.nightly-rust-clean-build</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/natemccoy/.claude/scripts/nightly/nightly-rust-clean-build.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/nightly-rust-clean-build-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/nightly-rust-clean-build-stderr.log</string>
</dict>
</plist>
