# scripts/nightly/com.natemccoy.nightly-rust-clean-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.natemccoy.nightly-rust-clean-build</string>
    <!-- StartCalendarInterval is unusable on this machine: UserEventAgent
         caches a stale Europe/Madrid TZ from prior travel and applies a
         6-hour offset to the scheduled hour. We trigger via StartInterval
         instead (raw seconds, no TZ involvement) and let the wrapper script
         gate on `date +%H`, which correctly reads /etc/localtime. The
         wrapper records the date of the most recent run so that only the
         first eligible fire of the day actually executes the build. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/natemccoy/.claude/scripts/nightly/nightly-trigger.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/nightly-rust-clean-build-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/nightly-rust-clean-build-stderr.log</string>
</dict>
</plist>
